# bot.conf
[program:moneyobmenbot]
directory=/home/ubuntu/moneyobmenbot
command=python3 app.py

autostart=true
autorestart=true
enviroment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/moneyobmenbot/logfile_err.log
stdout_logfile=/home/ubuntu/moneyobmenbot/logfile.log